# widget/helix-chat-api.service
[Unit]
Description=Helix Constitutional Chat API
After=network.target

[Service]
Type=simple
User=helix2
WorkingDirectory=/home/helix2/helix-adapter/widget
ExecStart=/home/helix2/.hermes/hermes-agent/venv/bin/python3 -m uvicorn api:app --host 127.0.0.1 --port 8001
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
